# deploy/parallax-tracker.service
[Unit]
Description=Monocular Parallax Head Tracker
After=network.target sound.target
Wants=network.target

[Service]
Type=simple
User=%i
Group=video
SupplementaryGroups=input
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/local/bin/parallax-daemon -c /etc/parallax/default.yaml
Restart=on-failure
RestartSec=2
Nice=-5
LimitMEMLOCK=infinity

[Install]
WantedBy=multi-user.target
